# deployment/soundcomparisons-backup.timer
[Unit]
Description=Run Sound-Comparisons backup script every 4 hours

[Timer]
OnBootSec=4h
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
